# packages/monitorat/monitorat-0.11.1.tar.gz/monitorat-0.11.1/systemd/monitor@pip.service
[Unit]
Description=%i Monitor
After=network.target

[Service]
User=__user__
Group=__group__
Environment="PATH=/home/__user__/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/env monitorat server --host 0.0.0.0 --port 6161
Restart=on-failure

[Install]
WantedBy=multi-user.target
