# packages/sysdi/sysdi-0.20260116.2.tar.gz/sysdi-0.20260116.2/src/sysdi_tests/systemd-units/aliases.timer
[Unit]
Description=Harvest Sync

[Timer]
AccuracySec=60s
OnActiveSec=5s
OnUnitActiveSec=5m
OnUnitInactiveSec=30s

[Install]
WantedBy=timers.target
