# packages/sysdi/sysdi-0.20260116.2.tar.gz/sysdi-0.20260116.2/src/sysdi_tests/systemd-units/daily.timer
[Unit]
Description=Daily Example

[Timer]
AccuracySec=60s
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3h
FixedRandomDelay=true

[Install]
WantedBy=timers.target
