# packages/sysdi/sysdi-0.20260116.2.tar.gz/sysdi-0.20260116.2/src/sysdi_tests/systemd-units/defaults.service
[Unit]
Description=check defaults

[Service]
Type=oneshot
ExecStart=/bin/picard foo bar
