# packages/sysdi/sysdi-0.20260116.2.tar.gz/sysdi-0.20260116.2/src/sysdi_tests/systemd-units/exec-wrap.service
[Unit]
Description=Check exec wrap

[Service]
Type=oneshot
ExecStart=/bin/picard
ExecStartPre=/usr/bin/exec-pre
ExecStopPost=/usr/bin/env bash -ec '\
    if [ "$SERVICE_RESULT" = success ]; then\
        /usr/bin/exec-post-ok;\
    else\
        /usr/bin/exec-post-fail;\
    fi'
