# packages/sysdi/sysdi-0.20260116.2.tar.gz/sysdi-0.20260116.2/src/sysdi_tests/systemd-units/extra-config.timer
[Unit]
Description=Extra Config

[Timer]
AccuracySec=60s
OnUnitInactiveSec=10

[Install]
WantedBy=timers.target
