# packages/sysdi/sysdi-0.20260116.2.tar.gz/sysdi-0.20260116.2/src/sysdi_tests/systemd-units/extra-config.service
[Unit]
Description=Extra Config

[Service]
Type=oneshot
ExecStart=/bin/picard
WorkingDirectory=foo
