# packages/wyoming-mlx-whisper/wyoming_mlx_whisper-1.4.0.tar.gz/wyoming_mlx_whisper-1.4.0/scripts/com.wyoming_mlx_whisper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wyoming_mlx_whisper</string>
  <key>ProgramArguments</key>
  <array>
    <string><UV-PATH></string>
    <string>tool</string>
    <string>run</string>
    <string>--python</string>
    <string>3.12</string>
    <string>wyoming-mlx-whisper</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string><HOME-DIR></string>
  <key>StandardOutPath</key>
  <string><LOG-DIR>/stdout.log</string>
  <key>StandardErrorPath</key>
  <string><LOG-DIR>/stderr.log</string>
</dict>
</plist>
