# packages/tetrascience-cli/tetrascience_cli-1.7.1-py3-none-any.whl/ts_cli/commands/init_cmd/templates/data-app/streamlit/4.3/conf/supervisord-local.conf
[supervisord]
nodaemon=true
user=root
pidfile=/tmp/supervisord.pid

[program:app]
command=python -m streamlit run app.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=30
stopsignal=TERM
killasgroup=true
stopasgroup=true
