# packages/valet-agent/valet_agent-2026.1.18.tar.gz/valet_agent-2026.1.18/valet-agent.service
[Unit]
Description=DeviceValet Linux Agent
Documentation=https://github.com/supported-systems/devicevalet
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/valet-agent
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/var/lib/valet-agent /var/log

# Environment
EnvironmentFile=-/etc/valet-agent/agent.env

[Install]
WantedBy=multi-user.target
