# packages/monitorat/monitorat-0.11-py3-none-any.whl/monitorat/systemd/monitor@source.service
[Unit]
Description=%i Monitor
After=network.target

[Service]
WorkingDirectory=__project__/monitorat
ExecStart=__project__/monitorat/.venv/bin/monitorat server --host 0.0.0.0 --port __port__
User=__user__
Group=__group__
Environment="PATH=__project__/monitorat/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Restart=on-failure

[Install]
WantedBy=multi-user.target
